# mdns-reflector/s6-rc.d/svc-mdns-reflector/run
#!/usr/bin/env sh
su irrexplorer -c 'mdns-reflector'
