# fullbogons/crontab/fullbogons-cron
0 0 */4 * * * * su bird -s /usr/bin/python3 fullbogons.py && birdc configure
