# irrd/s6-rc.d/init-irrexplorer/up
#!/command/execlineb -P
s6-setuidgid irrexplorer
/opt/irrexplorer/init-irrexplorer.sh
