# irrd/s6-rc.d/svc-irrd/run
#!/usr/bin/env sh
# irrd needs to start as root
exec /opt/irrd/bin/irrd --config /opt/irrd/irrd.yaml --foreground
