# irrd/s6-rc.d/svc-irrexplorer/run
#!/usr/bin/env sh
cd /opt/irrexplorer/ \
&& su irrexplorer -c '/usr/bin/poetry run python3 svc-irrexplorer.py'
