# growatt-ev/crontab
*/5 * * * * python3 /opt/growatt-ev/main.py
