# verlihub/s6-rc.d/svc-verlihub/run
#!/command/execlineb -P
with-contenv
s6-setuidgid verlihub
envfile /opt/verlihub/verlihub.env /usr/local/bin/verlihub
