# irrd/s6-rc.d/svc-supercronic/run
#!/usr/bin/env sh
su irrexplorer -c '/usr/local/bin/supercronic /opt/irrexplorer/cron/import-data'
