# irrd/s6-rc.d/init-irrd/up
#!/command/execlineb -P
s6-setuidgid irrd
pipx run /opt/irrd/init-irrd.sh
