# speedtest/s6-rc.d/init-speedtest/up
#!/command/execlineb -P
/opt/speedtest/speedtest.sh
