# speedtest/s6-rc.d/init-speedtest/up
#!/command/execlineb -P
/bin/sh -c /opt/speedtest/speedtest --accept-license --progress=no
